# t123.development.uwsgi.conf
#Project t123 uwsgi upstart job


description "Project t123 uwsgi upstart job"

start on runlevel [2345]
stop on runlevel [!2345]

env VE_HOME=/home/ilian/workspace/fabtests2/t123
env PROJECT_HOME=/home/ilian/workspace/fabtests2/t123/src/t123
env SOCKET=127.0.0.1:8002

env LOGTO=/var/log/nginx/t123_uwsgi.log
env RUN_PROCESSES=1
env DJANGO_SETTINGS_MODULE='t123.settings.development'
umask 007

# The default of 5 seconds is too low for mysql which needs to flush buffers
kill timeout 5

respawn

exec $VE_HOME/bin/uwsgi \
    --env DJANGO_SETTINGS_MODULE=$DJANGO_SETTINGS_MODULE \
    --processes $RUN_PROCESSES \
    --module=wsgi:application \
    --master \
    --socket $SOCKET \
    --uid=1000 --gid=2000 \
    --harakiri=20 \
    --limit-as=128 \
    --max-requests=5000 \
    --vacuum \
    --home=$VE_HOME \
    --pythonpath=$PROJECT_HOME \
    --logto $LOGTO >> /tmp/t123_uwsgi.log